# supervisor/serv.conf
[program:celery]
directory = /code
command = celery -A project worker -l warning
