# supervisor/serv.conf
[program:webserver]
directory = /opt/app
command = gunicorn project.wsgi:application --bind 0.0.0.0:8000

[program:celery]
directory = /opt/app
command = celery -A project worker -l warning
